# .service
[Unit]
Description=Gunicorn instance to serve web-api
After=network.target

[Service]
User=kali
Group=www-data
WorkingDirectory=/home/kali/web-api
Environment="PATH=/home/kali/web-api/env/bin"
ExecStart=/home/kali/web-api/env/bin/gunicorn --workers 3 --bind unix:web-api.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target